# data-collector-webservice/deploy/systemd.service
[Unit]
Description=Data Collector Web Service
After=network.target

[Service]
Type=simple
User=appuser
Group=appuser
WorkingDirectory=/opt/data-collector
EnvironmentFile=/etc/data-collector.env
ExecStart=/opt/data-collector/.venv/bin/gunicorn app.main:app \
    -k uvicorn.workers.UvicornWorker \
    --bind 127.0.0.1:8000 \
    --workers 4 \
    --timeout 60 \
    --access-logfile - \
    --error-logfile -
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/data-collector /nfs/airflow-staging

[Install]
WantedBy=multi-user.target
